# debian/doubao-input.service
[Unit]
Description=豆包语音输入法 (Doubao Input)
Documentation=https://github.com/wurong98/doubao-input-for-linux
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/doubao-input
Restart=on-failure
RestartSec=3

# Wayland / X11 都需要
Environment=DISPLAY=:0
Environment=GDK_BACKEND=wayland,x11
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target